# deploy/systemd/zerogex-web-trial-reminders.service
[Unit]
Description=ZeroGEX web ~48h-before-trial-end reminder emails
Documentation=https://github.com/zerogex-options/zerogex-web
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/zerogex-web

# `make trial-reminders` shells through nvm to get the correct Node 22
# binary (the script is a .mts that needs --experimental-strip-types). The
# nvm bash hook needs HOME to find ~/.nvm.
Environment="HOME=/home/ubuntu"
Environment="PATH=/home/ubuntu/.nvm/versions/node/v22/bin:/usr/local/bin:/usr/bin:/bin"

ExecStart=/usr/bin/make -C /home/ubuntu/zerogex-web trial-reminders YES=1

# Resend + Stripe are the only outbound legs; cap so a hung HTTPS connect
# can't pile up overlapping timer fires.
TimeoutStartSec=300

# Security hardening (mirrors the auth-backup unit).
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
# RW exception: the auth DB dir, because UPDATE users SET
# trial_reminder_email_sent_at and the audit_events INSERT both write to
# the same SQLite file as the live PM2 writer. Update this if you change
# AUTH_DB_PATH in frontend/.env.local.
ReadWritePaths=/home/ubuntu/zerogex-web/frontend/data

# Logging -> journald (view with: journalctl -u zerogex-web-trial-reminders)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=zerogex-web-trial-reminders
