# deploy/systemd/zerogex-web-winback.service
[Unit]
Description=ZeroGEX web ~1-month-after-churn win-back review digest
Documentation=https://github.com/zerogex-options/zerogex-web
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/zerogex-web

# `make winback` shells through nvm to get the correct Node 22 binary (the
# script is a .mts that needs --experimental-strip-types). The Makefile target
# sources $HOME/.nvm/nvm.sh internally, so all we need here is HOME pointing at
# the user's home — the system PATH suffices for /usr/bin/make + /usr/bin/bash;
# node lands on PATH via nvm.sh at runtime.
Environment="HOME=/home/ubuntu"
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

# REVIEW MODE, not auto-send: the weekly timer emails the founder the list of
# eligible churners + the rendered draft and sends NOTHING to users. The founder
# approves a batch by running `make winback YES=1` after reviewing. To make the
# weekly job auto-send instead, change this to `winback YES=1` (and note the
# job then writes the win-back latch, so ReadWritePaths below still applies).
# The digest recipient comes from WINBACK_DIGEST_TO in frontend/.env.local.
ExecStart=/usr/bin/make -C /home/ubuntu/zerogex-web winback DIGEST=1

# Resend + SQLite are the only outbound legs; cap so a hung HTTPS connect can't
# pile up overlapping timer fires.
TimeoutStartSec=300

# Security hardening (mirrors the verified-never-paid unit).
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
# RW exception: the auth DB dir. The digest run itself is read-only, but the
# approve step (`make winback YES=1`) writes UPDATE users SET
# winback_email_sent_at + the audit_events INSERT to the same SQLite file as the
# live PM2 writer, and shares this unit's sandbox if triggered via systemctl.
# Update this if you change AUTH_DB_PATH in frontend/.env.local. Deploy step 092
# also writes a drop-in that mirrors ReadWritePaths onto whatever you set there.
ReadWritePaths=/home/ubuntu/zerogex-web/frontend/data

# Logging -> journald (view with: journalctl -u zerogex-web-winback)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=zerogex-web-winback
